# config/supervisord.conf
; supervisord file for docker-compose webolith and macondo
; (and maybe our go websocket service after this!)


[program:dc_web]
directory = /home/ubuntu/webolith
command = docker-compose -f dc-prod.yml up
user=ubuntu
autostart=true
autorestart=true
priority=999

[program:macondo]
directory = /home/ubuntu/macondo
command = docker-compose -f docker-compose.yml -f dc-production.yaml up
user=ubuntu
autostart=true
autorestart=true
priority=991

[unix_http_server]
file=/tmp/supervisor.sock   ; path to your socket file

[supervisord]
logfile=/home/ubuntu/supervisord.log ; supervisord log file
logfile_maxbytes=50MB       ; maximum size of logfile before rotation
logfile_backups=10          ; number of backed up logfiles
loglevel=debug              ; info, debug, warn, trace
pidfile=/home/ubuntu/supervisord.pid ; pidfile location
nodaemon=false              ; run supervisord as a daemon
minfds=1024                 ; number of startup file descriptors
minprocs=200                ; number of process descriptors
user=root                   ; default user

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use unix:// schem for a unix sockets.

