# node/supervisord.conf
[program:mynode]
directory = /home/ubuntu/Webolith/node
command = bash -c "exec /usr/local/bin/node /home/ubuntu/Webolith/node/sockets.js"
priority=2
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 600
environment = NODE_ENV=production

[unix_http_server]
file=/tmp/supervisor.sock   ; path to your socket file

[supervisord]
logfile=/home/ubuntu/node_supervisord.log ; supervisord log file
logfile_maxbytes=50MB       ; maximum size of logfile before rotation
logfile_backups=10          ; number of backed up logfiles
loglevel=debug              ; info, debug, warn, trace
pidfile=/home/ubuntu/node_supervisord.pid ; pidfile location
nodaemon=false              ; run supervisord as a daemon
minfds=1024                 ; number of startup file descriptors
minprocs=200                ; number of process descriptors
user=root                   ; default user

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/node_supervisor.sock ; use unix:// schem for a unix sockets.

